# system/systemd/reservechain-update.service
[Unit]
Description=ReserveChain Auto Update

[Service]
Type=oneshot
User=reservechain
ExecStart=/opt/reservechain/scripts/update.sh
