# system/systemd/reservechain-node.service
[Unit]
Description=ReserveChain Node (Main Server + Seed)
After=network-online.target
Wants=network-online.target

[Service]
User=reservechain
WorkingDirectory=/opt/reservechain
ExecStart=/opt/reservechain/bin/reservechain-node --data-dir=/opt/reservechain/runtime/data --http=127.0.0.1:8080
Restart=always
RestartSec=3
LimitNOFILE=65536
StandardOutput=append:/opt/reservechain/runtime/logs/node.log
StandardError=append:/opt/reservechain/runtime/logs/node.err

[Install]
WantedBy=multi-user.target
